# overlays/opsec/etc/skel/.config/autostart/pipewire-pulse.desktop
[Desktop Entry]
Type=Application
Name=PipeWire PulseAudio
Exec=/usr/bin/pipewire-pulse
X-KDE-autostart-phase=2
