# overlays/opsec/etc/skel/.config/autostart/pipewire.desktop
[Desktop Entry]
Type=Application
Name=PipeWire
Exec=/usr/bin/pipewire
X-KDE-autostart-phase=1
