# overlays/opsec/etc/skel/.config/autostart/wireplumber.desktop
[Desktop Entry]
Type=Application
Name=WirePlumber
Exec=/usr/bin/wireplumber
X-KDE-autostart-phase=1
